# mysql/templates/mysql-down.conf.erb
start on stopping mysql

task

exec mysqldump  --password=<%= dbpass %> --user=<%= dbuser %> --opt <%= dbname %> > <%= dbdumppath %>/<%= dbname %>_dev.sql
